# chocolatey/private/Vagrant.EPS.Marketplace/tools/modules/elasticmq/templates/elasticmq.upstart.erb
description "elasticmq server <%= port %>"

start on runlevel [23]
stop on shutdown

exec sudo -u elasticmq /opt/elasticmq/bin/run.sh

respawn
